# batterang.timer
[Unit]
Description=Record battery healths every day

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=basic.target
