# batterang.service
[Unit]
Description=Record battery healths

[Service]
Type=simple
ExecStart=%h/.local/bin/batterang